# cron/crontab
# Jobber pipeline schedule - once daily at a non-round time to look less automated
23 8 * * * cd /app && sleep $((RANDOM \% 300)) && /usr/local/bin/python main.py >> /var/log/job-tracker/pipeline.log 2>&1
# empty line required by cron
