# cron/crontab
# Jobber pipeline schedule - runs at 7AM, 12PM, and 6PM daily
0 7  * * * cd /app && /usr/local/bin/python main.py >> /var/log/job-tracker/pipeline.log 2>&1
0 12 * * * cd /app && /usr/local/bin/python main.py >> /var/log/job-tracker/pipeline.log 2>&1
0 18 * * * cd /app && /usr/local/bin/python main.py >> /var/log/job-tracker/pipeline.log 2>&1
# empty line required by cron
